# crontabs/crontab
* * * * * root echo "Exemplo de Comando" >> /var/log/cron.log 2>&1
